# docker/mail/rootfs/etc/service/dovecot/run
#!/bin/sh
PID=$(pgrep dovecot)
if [ -z $PID ]
then
    exec /usr/sbin/dovecot -c /etc/mail/dovecot/dovecot.conf
fi
