# docker/mail/rootfs/etc/service/postfix/run
#!/bin/sh
PID=$(pgrep master)
if [ -z $PID ]
then
    exec /usr/sbin/postfix start -c /etc/mail/postfix
fi