# docker/mail/etc/service/rsyslog/run
#!/bin/sh
PID=$(pgrep rsyslog)
if [ -z "${PID}" ]
then
    exec /usr/sbin/rsyslogd -f /etc/mail/rsyslogd/rsyslog.conf
fi